# var/spool/cron/crontabs/root
# Delete offline messages that are too large
* * * * * bash /home/git/xmpp.is/scripts/delete-prosody-offline-messages.sh

# Update connection stats: https://transparency.xmpp.is/connection-stats.txt
* * * * * bash /home/git/xmpp.is/scripts/connection-stats.sh

# Force permissions
@hourly bash /home/git/xmpp.is/scripts/force-owner-and-group.sh

# Hourly GoAccess stats update
@hourly /usr/bin/goaccess /var/log/hiawatha/xmpp.is-access.log -o /var/www/stats.xmpp.is/index.html

# Hourly Borg Backup
@hourly bash /root/scripts/borgbackup.sh
