# var/spool/cron/crontabs/user
MAILTO=cron@lunar.systems

# Script to check if Tor HS v2/v3 is reachable
*/5 * * * * bash /home/user/git/xmpp.is/scripts/tor-hs-check.sh

# Renew certificate -- https://transparency.xmpp.is/fingerprints.txt
@weekly /home/user/git/dehydrated/dehydrated --cron --alias all --hook /home/user/git/dehydrated/dehydrated.default.sh; bash /home/user/git/xmpp.is/scripts/pre-le-renew-hook.sh >> /home/user/git/dehydrated/renew.log
