# var/spool/cron/crontabs/user
MAILTO=admin@lunar.systems

# Script to check if Tor HS v2/v3 is reachable
*/5 * * * * bash /home/user/git/xmpp.is/scripts/tor-hs-check.sh
